# config/systemd.service
[Unit]
Description=LLM Gateway - OpenAI-compatible proxy for on-demand model serving
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=%LLM_USER%
Group=%LLM_USER%
Environment=HOME=%LLM_HOME%
Environment=PATH=%LLM_PATH%
WorkingDirectory=/etc/llm-gateway
ExecStart=/usr/local/bin/llm-gateway

# --- GPU & Memory Permissions ---
SupplementaryGroups=render video
LimitMEMLOCK=infinity

# AMD/Intel DRI access (Radeon 780M iGPU, Intel iGPU)
DeviceAllow=/dev/dri/renderD128 rw
DeviceAllow=/dev/dri/renderD129 rw
DeviceAllow=/dev/dri/card0 rw
DeviceAllow=/dev/dri/card1 rw
# --------------------------------

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
